# 40-trigger-monitord-hook.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/systemctl start --no-block monitord-hook.service"
